# linux.lnk/.config/systemd/user/darkman.service.d/override.conf
[Service]
# DISPLAY, I3SOCK, XDG_CURRENT_DESKTOP are imported by i3 config via
# systemctl --user import-environment. We only need to prepend our
# custom xdg data dir. %h expands to the user's home directory.
Environment="XDG_DATA_DIRS=%h/.config/xdg:/usr/local/share:/usr/share"
