# config/systemd/user/darkman.service.d/override.conf
[Service]
Environment="XDG_DATA_DIRS=$HOME/.config/xdg:$XDG_DATA_DIRS"
